# debian/dynamic-ip.cron.d
*/15 * * * * root /usr/bin/dynamicIp.py
